# set-gpu-fans.service
[Unit]
Description=Cool GPU

[Service]
ExecStart=/opt/set-gpu-fans/cool_gpu
Restart=always

[Install]
Alias=set-gpu-fans.service

